# backend/supervisord.conf
[supervisord]

[program:websockets]
command = python3 ./websocket.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 1
autorestart = true